# ssh-proxy/ssh-proxy.service
[Unit]
Description=SSH WebSocket Proxy
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ssh-proxy
ExecStart=/usr/bin/node /opt/ssh-proxy/server.mjs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
